# systemd/nemotron-llm.service
[Unit]
Description=Nemotron-3-Nano LLM (vLLM)
After=network.target
Wants=network.target

[Service]
Type=simple
User=phil
Group=phil
WorkingDirectory=/home/phil/telephony-stack
Environment="PATH=/home/phil/telephony-stack-env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VLLM_WORKER_MULTIPROC_METHOD=spawn"
Environment="TRITON_PTXAS_PATH=/usr/local/cuda-13.0/bin/ptxas"
Environment="HF_HOME=/tmp/hf_cache"
Environment="CUDA_HOME=/usr/local/cuda-13.0"
Environment="PYTHONUNBUFFERED=1"
Environment="VLLM_ATTENTION_BACKEND=FLASHINFER"

ExecStart=/home/phil/telephony-stack-env/bin/python -m vllm.entrypoints.openai.api_server \
    --model /home/phil/telephony-stack/models/llm/nemotron-3-nano-30b-nvfp4 \
    --quantization modelopt_fp4 \
    --kv-cache-dtype fp8 \
    --gpu-memory-utilization 0.2 \
    --max-model-len 32768 \
    --enforce-eager \
    --trust-remote-code \
    --chat-template /home/phil/telephony-stack/telephony_template.jinja \
    --port 8000

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nemotron-llm

[Install]
WantedBy=multi-user.target
