# systemd/voxtral-asr.service
[Unit]
Description=Voxtral-Mini-4B-Realtime ASR (vLLM)
After=network.target nemotron-llm.service
Wants=network.target

[Service]
Type=simple
User=phil
Group=phil
WorkingDirectory=/home/phil/telephony-stack
Environment="PATH=/home/phil/telephony-stack-env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VLLM_WORKER_MULTIPROC_METHOD=spawn"
Environment="TRITON_PTXAS_PATH=/usr/local/cuda-13.0/bin/ptxas"
Environment="HF_HOME=/tmp/hf_cache"
Environment="CUDA_HOME=/usr/local/cuda-13.0"
Environment="PYTHONUNBUFFERED=1"

ExecStart=/home/phil/telephony-stack-env/bin/python -m vllm.entrypoints.openai.api_server \
    --model /home/phil/telephony-stack/models/asr/voxtral-mini-4b-realtime \
    --dtype bfloat16 \
    --max-model-len 8192 \
    --gpu-memory-utilization 0.1 \
    --enforce-eager \
    --trust-remote-code \
    --port 8001

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voxtral-asr

[Install]
WantedBy=multi-user.target
